# deploy/anyadeals-research.timer
[Unit]
Description=anyadeals AI research parsing timer

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
